# apache/www/memcached.conf
[program:memcached]
command=service memcached start
