# apache/www/php5.6-fpm.conf
[program:php5.6-fpm]
command=service php5.6-fpm start